# docker/systemd/system/mnt-nas-pve-volumes.automount
[Unit]
Description=Automount for enihsyou-NAS
ConditionPathExists=/mnt/nas/pve/volumes

[Automount]
Where=/mnt/nas/pve/volumes
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
