# network/docker/systemd/mount/mnt-nas-pve-volumes.automount
[Automount]
Where=/mnt/nas/pve/volumes
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
